# breachforge/intel/data/techniques.json
{
  "timestamp": "2026-05-29T20:05:18.435Z",
  "techniques": [
    {
      "id": "T1055.011",
      "name": "Extra Window Memory Injection",
      "description": "Adversaries may inject malicious code into process via Extra Window Memory (EWM) in order to evade p",
      "created": "2020-01-14T17:18:32.126Z"
    },
    {
      "id": "T1053.005",
      "name": "Scheduled Task",
      "description": "Adversaries may abuse the Windows Task Scheduler to perform task scheduling for initial or recurring",
      "created": "2019-11-27T14:58:00.429Z"
    },
    {
      "id": "T1205.002",
      "name": "Socket Filters",
      "description": "Adversaries may attach filters to a network socket to monitor then activate backdoors used for persi",
      "created": "2022-09-30T21:18:41.930Z"
    },
    {
      "id": "T1066",
      "name": "Indicator Removal from Tools",
      "description": "If a malicious tool is detected and quarantined or otherwise curtailed, an adversary may be able to ",
      "created": "2017-05-31T21:30:54.176Z"
    },
    {
      "id": "T1560.001",
      "name": "Archive via Utility",
      "description": "Adversaries may use utilities to compress and/or encrypt collected data prior to exfiltration. Many ",
      "created": "2020-02-20T21:01:25.428Z"
    },
    {
      "id": "T1021.005",
      "name": "VNC",
      "description": "Adversaries may use [Valid Accounts](https://attack.mitre.org/techniques/T1078) to remotely control ",
      "created": "2020-02-11T18:28:44.950Z"
    },
    {
      "id": "T1047",
      "name": "Windows Management Instrumentation",
      "description": "Adversaries may abuse Windows Management Instrumentation (WMI) to execute malicious commands and pay",
      "created": "2017-05-31T21:30:44.329Z"
    },
    {
      "id": "T1687",
      "name": "Exploitation for Defense Impairment",
      "description": "Adversaries may exploit vulnerabilities in security software, infrastructure, or defensive component",
      "created": "2026-04-14T22:53:27.621Z"
    },
    {
      "id": "T1156",
      "name": "Malicious Shell Modification",
      "description": "Adversaries may establish persistence through executing malicious commands triggered by a user’s she",
      "created": "2017-12-14T16:46:06.044Z"
    },
    {
      "id": "T1113",
      "name": "Screen Capture",
      "description": "Adversaries may attempt to take screen captures of the desktop to gather information over the course",
      "created": "2017-05-31T21:31:25.060Z"
    },
    {
      "id": "T1027.011",
      "name": "Fileless Storage",
      "description": "Adversaries may store data in \"fileless\" formats to conceal malicious activity from defenses. Filele",
      "created": "2023-03-23T19:55:25.546Z"
    },
    {
      "id": "T1067",
      "name": "Bootkit",
      "description": "A bootkit is a malware variant that modifies the boot sectors of a hard drive, including the Master ",
      "created": "2017-05-31T21:30:54.661Z"
    },
    {
      "id": "T1037",
      "name": "Boot or Logon Initialization Scripts",
      "description": "Adversaries may use scripts automatically executed at boot or logon initialization to establish pers",
      "created": "2017-05-31T21:30:38.910Z"
    },
    {
      "id": "T1557",
      "name": "Adversary-in-the-Middle",
      "description": "Adversaries may attempt to position themselves between two or more networked devices using an advers",
      "created": "2020-02-11T19:07:12.114Z"
    },
    {
      "id": "T1033",
      "name": "System Owner/User Discovery",
      "description": "Adversaries may attempt to identify the primary user, currently logged in user, set of users that co",
      "created": "2017-05-31T21:30:35.733Z"
    }
  ],
  "total": 858
}